# squid-proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:squid4]
command=/apps/squid/sbin/squid -NsY -f /apps/squid.conf.intercept-cache
stdout_events_enabled=true
stderr_events_enabled=true


